# .config/systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -b -w 10 -W Battery level is low -c 5 -C Battery level is critical -a Battery
